# files/systemd/mouseless.service
[Unit]
Description=mouseless

[Service]
ExecStart=/usr/local/bin/mouseless --debug --config /home/anthony/.config/mouseless/config.yaml

[Install]
WantedBy=multi-user.target
